# scripts/edsys-share/systemd/edsys-share-gdrive-prune.service
[Unit]
Description=Prune expired EdSys Share Drive recovery folders to Trash
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/EdSys-Share /mnt/ai-store

[Service]
Type=oneshot
UMask=0077
Environment=EDSYS_SHARE_REQUIRE_WRITABLE=false
ExecStart=/usr/local/sbin/edsys-share-gdrive-prune
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/etc/edsys-share
ReadWritePaths=/etc/edsys-backup/rclone.conf /var/lib/edsys-share /var/cache/edsys-share /srv/edsys-backup/reports/edsys-share /run
NoNewPrivileges=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
