# scripts/security/edsys-ssh-guard.service
[Unit]
Description=Restrict 9950x SSH to intended LAN and Tailnet controllers
# ssh.socket normally starts before sockets.target.  Run this firewall guard
# early enough that requiring it from ssh.socket cannot create a boot ordering
# cycle through basic.target -> sockets.target.
DefaultDependencies=no
After=local-fs.target
Before=network-pre.target ssh.socket ssh.service
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/edsys-ssh-guard
ExecReload=/usr/local/sbin/edsys-ssh-guard
ExecStop=-/usr/sbin/nft delete table inet edsys_ssh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
RequiredBy=ssh.socket ssh.service
