# scripts/ops/systemd/edsys-rag-golden-eval.service
[Unit]
Description=Run the private EdSys grounded-RAG golden-query quality gate
Documentation=file:/srv/edsys/edsys-infrastructure/scripts/ops/README.md
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/srv/edsys /mnt/ai-store
AssertPathIsMountPoint=/mnt/ai-store
AssertPathExists=/etc/edsys-secrets/edsys-rag-eval.env
AssertPathExists=/srv/edsys/EdSys-Master/tools/rag-eval/edsys-rag-golden-eval.py
AssertPathExists=/srv/edsys/EdSys-Master/data/rag-golden-queries.yml
AssertPathExists=/home/jeremy/code/edsys-ai-portal/.venv/bin/python
AssertPathExists=/home/jeremy/code/edsys-ai-portal/tools/rag_eval_runner.py

[Service]
Type=oneshot
User=jeremy
Group=jeremy
WorkingDirectory=/var/lib/edsys-rag-eval
UMask=0077
StateDirectory=edsys-rag-eval
StateDirectoryMode=0700
EnvironmentFile=/etc/edsys-secrets/edsys-rag-eval.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 /srv/edsys/EdSys-Master/tools/rag-eval/edsys-rag-golden-eval.py run --queries /srv/edsys/EdSys-Master/data/rag-golden-queries.yml --portal-repo /home/jeremy/code/edsys-ai-portal --portal-url http://192.168.50.50:3020 --judge heuristic --model edsys-grounded-local --results-dir /mnt/ai-store/rag/evals --state-dir /var/lib/edsys-rag-eval
TimeoutStartSec=45min
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
PrivateTmp=true
PrivateDevices=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
InaccessiblePaths=-/home/jeremy/code/edsys-ai-portal/.env -/home/jeremy/.codex
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/edsys-rag-eval /mnt/ai-store/rag/evals
