# scripts/backup/systemd/edsys-offsite-sync.service
[Unit]
Description=EdSys Google Drive offsite backup mirror
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/srv/edsys-backup/scripts/edsys-offsite-sync.sh
