# scripts/backup/systemd/edsys-backup.timer
[Unit]
Description=Run EdSys encrypted offsite backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
