# scripts/backup/systemd/edsys-backup-check.service
[Unit]
Description=EdSys restic repository health check
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/srv/edsys-backup/scripts/edsys-restic-check.sh
