# scripts/backup/systemd/edsys-restore-test.service
[Unit]
Description=EdSys monthly restore test
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/srv/edsys-backup/scripts/edsys-restore-test.sh
