# docker/edsys-courier/courier-smb-tailnet.socket.example
[Unit]
Description=EdSys Courier SMB listener on the Tailscale address
Documentation=man:systemd.socket(5)
After=tailscaled.service

[Socket]
# Replace this documentation-only value during deployment. Tailscale device IPs
# are stable, but this unit must be updated if the 9950x Tailnet IP changes.
ListenStream=100.64.0.1:445
BindToDevice=tailscale0
NoDelay=true

[Install]
WantedBy=sockets.target
